# apps/systemd/.config/systemd/user/battery-checker.timer
[Unit]
Description=battery checker timer

[Timer]
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
